# crontab
* * * * * php app/console solr:indexer:run
* * * * * php app/console locking:dbal:clean
